# settings/supervisor/conf.d/php-fpm.conf
[program:php-fpm]
command=/home/9zilla/.anyenv/envs/phpenv/versions/7.0.14/sbin/php-fpm --nodaemonize
autostart=true
autorestart=true
